# systemd/blockhost-gc.timer
[Unit]
Description=Daily BlockHost VM garbage collection
Documentation=man:blockhost-vm-gc(1)

[Timer]
# Run daily at 2 AM with up to 30 minutes of randomized delay
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
