# systemd/blockhost-gc.service
[Unit]
Description=BlockHost VM Garbage Collection
Documentation=man:blockhost-vm-gc(1)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/blockhost-vm-gc --execute
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/blockhost

[Install]
WantedBy=multi-user.target
